# templates/cleanup.service.tmpl
[Unit]
Description=Cleanup old MCP workspaces

[Service]
Type=oneshot
ExecStart=/usr/bin/find {{WORKSPACES_DIR}} -mindepth 1 -maxdepth 1 -type d -mtime +7 -print -exec rm -rf {} \;
